# debian/dbpagerd.upstart
# dbPager Server Upstart job configuration file
author "Dennis Prochko <wolfsoft@mail.ru>"
description "dbPager Server"

start on runlevel [2345]
stop on runlevel [016]

expect fork
respawn

script
	DBPAGER=/usr/bin/$UPSTART_JOB
	DBPAGER_OPTS=
	if [ -f /etc/default/$UPSTART_JOB ]; then
		. /etc/default/$UPSTART_JOB
	fi
	exec $DBPAGER $DBPAGER_OPTS -p /var/run/$UPSTART_JOB.pid
end script

post-stop script
	kill `cat /var/run/$UPSTART_JOB.pid`
	rm -f /var/run/$UPSTART_JOB.pid
end script
